# roles/prometheus-node-setup/nginx_exporter/templates/nginx-prometheus-exporter.socket.j2
[Unit]
Description=NGINX Prometheus Exporter Socket

[Socket]
ListenStream={{ NGINX_EXPORTER_SOCKET_PORT }}

[Install]
WantedBy=sockets.target
